# deploy/systemd/fap-web.service
[Unit]
Description=fap-web (Next.js standalone)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/fap-web
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
Environment=NEXT_PUBLIC_API_URL=https://api.fermatmind.com
Environment=NEXT_PUBLIC_RELEASE=web@2026-02-23+launch_gate
Environment=ENABLE_MBTI=true
Environment=ENABLE_BIG5_OCEAN=true
Environment=ENABLE_SDS_20=true
Environment=ENABLE_CLINICAL_COMBO_68=true
Environment=ROLLOUT_PERCENT_MBTI=100
Environment=ROLLOUT_PERCENT_BIG5_OCEAN=100
Environment=ROLLOUT_PERCENT_SDS_20=100
Environment=ROLLOUT_PERCENT_CLINICAL_COMBO_68=100
Environment=ENABLE_SDS_20_COMMERCE=true
Environment=ENABLE_CLINICAL_COMBO_68_COMMERCE=true
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
