# deploy/systemd/fap-web.service
# Aligned fallback/reference unit for hosts that explicitly install fap-web.service.
# Current Node1 production authority remains PM2 with /usr/bin/node on Node 24.x.
[Unit]
Description=fap-web (Next.js standalone reference unit)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/apps/fap-web/.next/standalone
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=FERMATMIND_DEPLOYED_REVISION_FILE=/opt/apps/fap-web/REVISION
# /usr/bin/node must resolve to Node 24.x before using this fallback/reference unit.
ExecStart=/usr/bin/node /opt/apps/fap-web/.next/standalone/server.js --hostname 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
