# docs/deploy/systemd-fap-web.service
# Reference systemd unit for hosts that explicitly install fap-web.service.
# Current Node1 production authority remains PM2, not this unit.
# /usr/bin/node must resolve to Node 24.x before using this fallback/reference file.

[Unit]
Description=fap-web (Next.js standalone reference unit)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/apps/fap-web/.next/standalone
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000
Environment=FERMATMIND_DEPLOYED_REVISION_FILE=/opt/apps/fap-web/REVISION
ExecStart=/usr/bin/node /opt/apps/fap-web/.next/standalone/server.js --hostname 127.0.0.1 --port 3000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
